# systemd/claude-watchdog.service
[Unit]
Description=Claude Watchdog - mode-aware forced restart

[Service]
Type=oneshot
EnvironmentFile=%h/.config/claude-restart/env
ExecStart=/bin/bash -c '\
  if [ "$CLAUDE_CONNECT" = "remote-control" ]; then \
    echo "claude-watchdog: skipped restart (remote-control mode has built-in reconnection)"; \
    exit 0; \
  fi; \
  echo "claude-watchdog: restarting claude service (mode=$CLAUDE_CONNECT)"; \
  systemctl --user restart claude.service'
