# systemd/claude.service
[Unit]
Description=Claude Code AI Assistant
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
ExecStart=%h/.local/bin/claude-wrapper --dangerously-skip-permissions
EnvironmentFile=%h/.config/claude-restart/env
WorkingDirectory=WORKING_DIR_PLACEHOLDER
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=default.target
